# crontab.cron
#Mins Hours Days Months Day of the week
#cron scripts need to be in path
00	20	*	*	*	rsync-backup.sh
00	15	*	*	*	yes | sudo pacman -Syu > ~/pacman.log
*	*	*	*	*	fcrontab -l > ~/Apps/stuff/crontab.cron
