# cron/crontab.cron
#Mins Hours Days Months Day of the week
#rsync-backup daily, cron scripts need to be in path
00	20	*	*	*	rsync-backup.sh
